# buildroot-v86/board/v86/rootfs_overlay/etc/inittab
# This is run first except when booting in single-user mode.
#
#::sysinit:/etc/boot.sh

# /bin/sh invocations on selected ttys
#
# Start an "respawn" shell on the console (whatever that may be)
::respawn:-/bin/sh
# Start an "askfirst" shell on /dev/tty2-4
tty2::askfirst:-/bin/sh
tty3::askfirst:-/bin/sh
tty4::askfirst:-/bin/sh

# /sbin/getty invocations for selected ttys
#
#tty4::respawn:/sbin/getty 38400 tty5
#tty5::respawn:/sbin/getty 38400 tty6

# Example of how to put a getty on a serial line (for a terminal)
#
ttyS0::respawn:-/bin/sh
#::respawn:/sbin/getty -L ttyS0 9600 vt100
#::respawn:/sbin/getty -L ttyS1 9600 vt100
#
# Example how to put a getty on a modem line.
#::respawn:/sbin/getty 57600 ttyS2

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot -f
#::ctrlaltdel:/sbin/swapoff
